# data/robopipeapi.service
[Unit]
Description=Start Robopipe API
After=robopipehostname.service
Requires=robopipehostname.service

[Service]
ExecStart=/bin/sh -c 'HOSTNAME=$(hostname) robopipe-api'
Restart=always
User=admin
Group=admin
Environment="ROBOPIPE_API_ENV=/etc/robopipe/.env"

[Install]
WantedBy=multi-user.target